# dist/smarthome-hw-ir.service
[Unit]
Description=Infrared receiver for Smarthome
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/usr/bin/smarthome-hw-ir/
ExecStart=/usr/bin/smarthome-hw-ir/smarthome-hw-ir
Restart=always

[Install]
WantedBy=multi-user.target
